# www/yopass/pkg-plist
bin/yopass
bin/yopass-server
etc/rc.d/yopass
%%WWWDIR%%/public/_redirects
%%WWWDIR%%/public/assets/browser-ponyfill-JACPZ65_.js
%%WWWDIR%%/public/assets/browser-ponyfill-JACPZ65_.js.map
%%WWWDIR%%/public/assets/index-B1qGElgd.css
%%WWWDIR%%/public/assets/index-eol_BJFf.js
%%WWWDIR%%/public/assets/index-eol_BJFf.js.map
%%WWWDIR%%/public/assets/roboto-cyrillic-300-normal-BoeCr7ev.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-300-normal-brFHmz5Z.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-400-normal-DAIM1_dR.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-400-normal-Dry59Hjn.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-500-normal-CsCirF4J.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-500-normal-hCeO1jFL.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-700-normal-CB1Rmiii.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-700-normal-CzEIZVQR.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-300-normal-CUPJdUZp.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-300-normal-CbEu6Uxw.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-400-normal-BxX1-eA_.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-400-normal-DzMWdK87.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-500-normal-B7rQpwPu.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-500-normal-ElvJfk8V.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-700-normal-CIu0AXX2.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-700-normal-D_fA0fHY.woff2
%%WWWDIR%%/public/assets/roboto-greek-300-normal-C_kb6SiH.woff
%%WWWDIR%%/public/assets/roboto-greek-300-normal-Cnub06j7.woff2
%%WWWDIR%%/public/assets/roboto-greek-400-normal-CAI06USH.woff
%%WWWDIR%%/public/assets/roboto-greek-400-normal-jFM2czAU.woff2
%%WWWDIR%%/public/assets/roboto-greek-500-normal-BJMS0heP.woff2
%%WWWDIR%%/public/assets/roboto-greek-500-normal-D8eQD5zT.woff
%%WWWDIR%%/public/assets/roboto-greek-700-normal-DCrt6r9F.woff
%%WWWDIR%%/public/assets/roboto-greek-700-normal-DpKAje7q.woff2
%%WWWDIR%%/public/assets/roboto-latin-300-normal-CztWkFGs.woff2
%%WWWDIR%%/public/assets/roboto-latin-300-normal-xb4kL3zn.woff
%%WWWDIR%%/public/assets/roboto-latin-400-normal-551zQQ7R.woff
%%WWWDIR%%/public/assets/roboto-latin-400-normal-CNwBRw8h.woff2
%%WWWDIR%%/public/assets/roboto-latin-500-normal-CkrA1NAy.woff2
%%WWWDIR%%/public/assets/roboto-latin-500-normal-_8jDuD7w.woff
%%WWWDIR%%/public/assets/roboto-latin-700-normal-CTLkNcF_.woff
%%WWWDIR%%/public/assets/roboto-latin-700-normal-DZr4b_KL.woff2
%%WWWDIR%%/public/assets/roboto-latin-ext-300-normal-BOz0oSmX.woff2
%%WWWDIR%%/public/assets/roboto-latin-ext-300-normal-ni96UL8R.woff
%%WWWDIR%%/public/assets/roboto-latin-ext-400-normal-ZYmyxeOy.woff2
%%WWWDIR%%/public/assets/roboto-latin-ext-400-normal-uRIBRJt5.woff
%%WWWDIR%%/public/assets/roboto-latin-ext-500-normal-C_ARlJGk.woff2
%%WWWDIR%%/public/assets/roboto-latin-ext-500-normal-eJ10kk0m.woff
%%WWWDIR%%/public/assets/roboto-latin-ext-700-normal-BNPgmEQS.woff2
%%WWWDIR%%/public/assets/roboto-latin-ext-700-normal-Cnx4FGpK.woff
%%WWWDIR%%/public/assets/roboto-math-300-normal-BMxwzQmU.woff2
%%WWWDIR%%/public/assets/roboto-math-300-normal-Bhb_hbXj.woff
%%WWWDIR%%/public/assets/roboto-math-400-normal-B3wgz80t.woff2
%%WWWDIR%%/public/assets/roboto-math-400-normal-DHrwdhE6.woff
%%WWWDIR%%/public/assets/roboto-math-500-normal-CFNaIMFC.woff2
%%WWWDIR%%/public/assets/roboto-math-500-normal-CetgDdIa.woff
%%WWWDIR%%/public/assets/roboto-math-700-normal-CpTCM92H.woff
%%WWWDIR%%/public/assets/roboto-math-700-normal-xbpggnJp.woff2
%%WWWDIR%%/public/assets/roboto-symbols-300-normal-BtHl5mYl.woff2
%%WWWDIR%%/public/assets/roboto-symbols-300-normal-VgGNgJlY.woff
%%WWWDIR%%/public/assets/roboto-symbols-400-normal-bG5rsNFs.woff
%%WWWDIR%%/public/assets/roboto-symbols-400-normal-fF1SLJBj.woff2
%%WWWDIR%%/public/assets/roboto-symbols-500-normal-BXFTxrNR.woff2
%%WWWDIR%%/public/assets/roboto-symbols-500-normal-toKUCDph.woff
%%WWWDIR%%/public/assets/roboto-symbols-700-normal-B2QKVW64.woff
%%WWWDIR%%/public/assets/roboto-symbols-700-normal-DKkQdRpM.woff2
%%WWWDIR%%/public/assets/roboto-vietnamese-300-normal-CNeuLW5X.woff2
%%WWWDIR%%/public/assets/roboto-vietnamese-300-normal-D46JPfnP.woff
%%WWWDIR%%/public/assets/roboto-vietnamese-400-normal-CDDxGrUb.woff2
%%WWWDIR%%/public/assets/roboto-vietnamese-400-normal-DgufTq8s.woff
%%WWWDIR%%/public/assets/roboto-vietnamese-500-normal-Dw5heWgq.woff
%%WWWDIR%%/public/assets/roboto-vietnamese-500-normal-HYpufUYk.woff2
%%WWWDIR%%/public/assets/roboto-vietnamese-700-normal-BFWtvCOj.woff2
%%WWWDIR%%/public/assets/roboto-vietnamese-700-normal-ChAl_rRV.woff
%%WWWDIR%%/public/favicon-16x16.png
%%WWWDIR%%/public/favicon-32x32.png
%%WWWDIR%%/public/favicon.ico
%%WWWDIR%%/public/favicon.png
%%WWWDIR%%/public/index.css
%%WWWDIR%%/public/index.html
%%WWWDIR%%/public/locales/en.json
%%WWWDIR%%/public/manifest.json
%%WWWDIR%%/public/safari-pinned-tab.svg
%%WWWDIR%%/public/yopass.svg
%%WWWDIR%%/yopass-server
@dir /var/run/yopass
