# www/yopass/pkg-plist
bin/yopass-server
etc/rc.d/yopass
%%WWWDIR%%/public/_redirects
%%WWWDIR%%/public/assets/index-ChvkVzSO.css
%%WWWDIR%%/public/assets/index-Dae3YNuo.js
%%WWWDIR%%/public/assets/index-Dae3YNuo.js.map
%%WWWDIR%%/public/assets/roboto-cyrillic-300-normal-DJfICpyc.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-300-normal-Dg7J0kAT.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-400-normal-BiRJyiea.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-400-normal-JN0iKxGs.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-500-normal-YnJLGrUm.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-500-normal-_hamcpv8.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-700-normal-BJaAVvFw.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-700-normal-jruQITdB.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-300-normal-BLLmCegk.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-300-normal-Chhwl1Jq.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-400-normal-D76n7Daw.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-400-normal-b0JluIOJ.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-500-normal-37WQE4S0.woff
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-500-normal-BJvL3D7h.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-700-normal-CyZgh00P.woff2
%%WWWDIR%%/public/assets/roboto-cyrillic-ext-700-normal-DXzexxfu.woff
%%WWWDIR%%/public/assets/roboto-greek-300-normal-Bx8edVml.woff2
%%WWWDIR%%/public/assets/roboto-greek-300-normal-D3gN5oZ1.woff
%%WWWDIR%%/public/assets/roboto-greek-400-normal-IIc_WWwF.woff
%%WWWDIR%%/public/assets/roboto-greek-400-normal-LPh2sqOm.woff2
%%WWWDIR%%/public/assets/roboto-greek-500-normal-Bg8BLohm.woff2
%%WWWDIR%%/public/assets/roboto-greek-500-normal-CdRewbqV.woff
%%WWWDIR%%/public/assets/roboto-greek-700-normal-1IZ-NEfb.woff
%%WWWDIR%%/public/assets/roboto-greek-700-normal-Bs05n1ZH.woff2
%%WWWDIR%%/public/assets/roboto-latin-300-normal-BZ6gvbSO.woff
%%WWWDIR%%/public/assets/roboto-latin-300-normal-BizgZZ3y.woff2
%%WWWDIR%%/public/assets/roboto-latin-400-normal-BVyCgWwA.woff
%%WWWDIR%%/public/assets/roboto-latin-400-normal-DXyFPIdK.woff2
%%WWWDIR%%/public/assets/roboto-latin-500-normal-C6iW8rdg.woff2
%%WWWDIR%%/public/assets/roboto-latin-500-normal-rpP1_v3s.woff
%%WWWDIR%%/public/assets/roboto-latin-700-normal-BWcFiwQV.woff
%%WWWDIR%%/public/assets/roboto-latin-700-normal-CbYYDfWS.woff2
%%WWWDIR%%/public/assets/roboto-latin-ext-300-normal-BzRVPTS2.woff2
%%WWWDIR%%/public/assets/roboto-latin-ext-300-normal-Djx841zm.woff
%%WWWDIR%%/public/assets/roboto-latin-ext-400-normal-BSFkPfbf.woff
%%WWWDIR%%/public/assets/roboto-latin-ext-400-normal-DgXbz5gU.woff2
%%WWWDIR%%/public/assets/roboto-latin-ext-500-normal-DvHxAkTn.woff
%%WWWDIR%%/public/assets/roboto-latin-ext-500-normal-OQJhyaXd.woff2
%%WWWDIR%%/public/assets/roboto-latin-ext-700-normal-Ba-CAIIA.woff
%%WWWDIR%%/public/assets/roboto-latin-ext-700-normal-DchBbzVz.woff2
%%WWWDIR%%/public/assets/roboto-vietnamese-300-normal-CAomnZLO.woff
%%WWWDIR%%/public/assets/roboto-vietnamese-300-normal-PZa9KE_J.woff2
%%WWWDIR%%/public/assets/roboto-vietnamese-400-normal-D5pJwT9g.woff
%%WWWDIR%%/public/assets/roboto-vietnamese-400-normal-DhTUfTw_.woff2
%%WWWDIR%%/public/assets/roboto-vietnamese-500-normal-LvuCHq7y.woff
%%WWWDIR%%/public/assets/roboto-vietnamese-500-normal-p0V0BAAE.woff2
%%WWWDIR%%/public/assets/roboto-vietnamese-700-normal-B4Nagvlm.woff
%%WWWDIR%%/public/assets/roboto-vietnamese-700-normal-CBbheh0s.woff2
%%WWWDIR%%/public/favicon-16x16.png
%%WWWDIR%%/public/favicon-32x32.png
%%WWWDIR%%/public/favicon.ico
%%WWWDIR%%/public/favicon.png
%%WWWDIR%%/public/index.html
%%WWWDIR%%/public/locales/en.json
%%WWWDIR%%/public/manifest.json
%%WWWDIR%%/public/safari-pinned-tab.svg
%%WWWDIR%%/public/yopass.svg
%%WWWDIR%%/yopass-server
@dir /var/run/yopass
